# supervisor.ini
[program:analyzer]
command=/usr/bin/python /opt/monitor.analyzer/monitor/agent_analyzer_v3.py
directory=/opt/monitor.analyzer/monitor

redirect_stderr=true
stdout_logfile=/var/log/supervisord/analyzer.log
stderr_logfile=/var/log/supervisord/analyzer_error.log

user=root
autostart=true
autorestart=true
priority=991

stopsignal=KILL
killasgroup=true
stopasgroup=true

[program:analyzer_check]
command=/usr/bin/python /opt/monitor.analyzer/monitor/agent_analyzer_check_v3.py
directory=/opt/monitor.analyzer/monitor

redirect_stderr=true
stdout_logfile=/var/log/supervisord/analyzer_check.log
stderr_logfile=/var/log/supervisord/analyzer_check_error.log

user=root
autostart=true
autorestart=true
priorit
